# scripts/com.freeyolo.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.freeyolo.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zzzhu/workspace/free-yolo/scripts/daily.sh</string>
    </array>

    <!-- Run every day at 09:00 local time. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the Mac was asleep/off at 09:00, run as soon as it wakes. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/zzzhu/workspace/free-yolo/data/daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zzzhu/workspace/free-yolo/data/daily.log</string>
</dict>
</plist>
